# sandbox/worker/s6/signal-bridge/finish
#!/command/with-contenv bash
# Runs after signal-bridge exits. Cleans up orphaned subshells/tails.
#
# Does NOT touch /state/worker/lease-held — that file is shared with any
# blue/green peer worker. A departing bridge's cleanup trap already clears
# the sentinel if this bridge owned it; if a peer just took over, removing
# the file here would wrongly signal "no one holds the lease."
pkill -9 -f /usr/local/bin/signal-bridge 2>/dev/null || true
pkill -9 -f "tail -F.*signal-daemon.log" 2>/dev/null || true
exit 0
