# sandbox/worker/s6/signal-bridge/run
#!/command/with-contenv bash
# Route Signal messages to local claude. Runs as alice.
# Waits for the daemon's HTTP port to come up (daemon container's health
# check usually gates this via depends_on, but belt-and-braces).
set -e
for _ in $(seq 1 60); do
    if curl -sfo /dev/null -X POST "${SIGNAL_API:-http://alice-daemon:8080}/api/v1/rpc" \
        -H 'Content-Type: application/json' \
        -d '{"jsonrpc":"2.0","method":"version","id":"ping"}' 2>/dev/null; then
        break
    fi
    sleep 1
done

export ALICE_CONFIG=/home/alice/.config/alice/alice.env
exec s6-setuidgid alice /usr/local/bin/signal-bridge
