# sandbox/worker/s6/alice-speaking/finish
#!/command/with-contenv bash
# Cleanup after alice-speaking exits. The Python daemon already drains
# cleanly on SIGTERM; this is just a belt-and-suspenders kill for any
# stray subprocesses the Agent SDK might have spawned (claude CLI).
pkill -9 -f 'python -m alice_speaking' 2>/dev/null || true
pkill -9 -f '/opt/alice-venv/bin/python' 2>/dev/null || true
exit 0
