# sandbox/daemon/s6/signal-daemon/run
#!/command/with-contenv bash
# signal-cli HTTP daemon. Runs as alice.
# stdout/stderr → /state/daemon/signal-daemon.log (a shared volume the worker
# container tails for inbound messages).
#
# When SIGNAL_ACCOUNT is unset (alice.env hasn't been wired for Signal),
# the service idles — the daemon container stays up so compose's
# service_healthy gate can be satisfied via the disabled-mode marker
# the healthcheck script writes.
set -e
set -a
# shellcheck disable=SC1091
[ -f /home/alice/.config/alice/alice.env ] && source /home/alice/.config/alice/alice.env
set +a

LOG_FILE=/state/daemon/signal-daemon.log
mkdir -p "$(dirname "$LOG_FILE")"
touch "$LOG_FILE"
chown alice:alice "$LOG_FILE"

if [ -z "${SIGNAL_ACCOUNT:-}" ]; then
    echo "[signal-daemon] SIGNAL_ACCOUNT not set — Signal disabled, idling" >> "$LOG_FILE"
    # Mark disabled so the healthcheck reports healthy without curl-ing
    # signal-cli (which isn't running).
    mkdir -p /state/daemon
    : > /state/daemon/signal-disabled
    exec sleep infinity
fi
rm -f /state/daemon/signal-disabled

# Bind 0.0.0.0:8080 so the worker container can reach via Docker network.
exec s6-setuidgid alice bash -c \
    'exec signal-cli -a "$1" -o json daemon --http 0.0.0.0:8080 >> "$2" 2>&1' \
    -- "$SIGNAL_ACCOUNT" "$LOG_FILE"
