# sandbox/s6/signal-bridge/run
#!/command/with-contenv bash
# Route Signal messages to local claude. Runs as alice.
# Waits briefly for signal-daemon's HTTP port before starting (s6 already
# ensures signal-daemon has been started, but the JVM-free native binary
# can still take a moment to bind :8080).
set -e
for _ in $(seq 1 60); do
    if curl -sfo /dev/null -X POST http://127.0.0.1:8080/api/v1/rpc \
        -H 'Content-Type: application/json' \
        -d '{"jsonrpc":"2.0","method":"listAccounts","id":"ping"}' 2>/dev/null \
        || [ $? -eq 22 ]; then
        break
    fi
    sleep 1
done

export ALICE_CONFIG=/home/alice/.config/alice/alice.env
exec s6-setuidgid alice /usr/local/bin/signal-bridge
