# sandbox/s6/signal-daemon/run
#!/command/with-contenv bash
# signal-cli HTTP daemon. Runs as alice.
# Output goes to /home/alice/.local/state/alice/signal-daemon.log — the
# bridge tails that file instead of systemd's journal (which doesn't exist
# in this container).
set -e
set -a
# shellcheck disable=SC1091
source /home/alice/.config/alice/alice.env
set +a
: "${SIGNAL_ACCOUNT:?SIGNAL_ACCOUNT must be set in alice.env}"

LOG_FILE=/home/alice/.local/state/alice/signal-daemon.log
mkdir -p "$(dirname "$LOG_FILE")"
touch "$LOG_FILE"
chown alice:alice "$LOG_FILE"

# Use exec with redirect so signal-cli's stdout/stderr land in the log file;
# s6-setuidgid drops privileges before exec.
exec s6-setuidgid alice bash -c 'exec signal-cli -a "$1" -o json daemon --http 127.0.0.1:8080 >> "$2" 2>&1' -- "$SIGNAL_ACCOUNT" "$LOG_FILE"
