# sandbox/worker/s6/alice-gh-watcher/run
#!/command/with-contenv bash
# alice-gh-watcher — timer loop that polls watched GitHub repos and drops
# one note per unseen PR/review/comment/check event into
# ~/alice-mind/inner/notes/ for thinking to drain.
#
# Singleton via `flock -n` so a slow poll can't double-fire.
#
# Cadence is read from alice.config.json `github_watcher.poll_interval_seconds`
# on each cycle (live reconfig, no restart needed). Default: 300s. Watcher
# is a no-op when no repos are configured, so it's safe to leave running.
set -e

LOCK=/state/worker/gh-watcher.lock
CONFIG_JSON=/home/alice/alice-mind/config/alice.config.json
DEFAULT_INTERVAL=300

mkdir -p "$(dirname "$LOCK")"
touch "$LOCK"
chown alice:alice "$LOCK" || true

# Stagger startup so we don't all wake up in lockstep when the worker boots.
sleep 30

while true; do
    interval="$DEFAULT_INTERVAL"
    if [ -f "$CONFIG_JSON" ] && command -v jq >/dev/null 2>&1; then
        v=$(jq -r '.github_watcher.poll_interval_seconds // empty' "$CONFIG_JSON" 2>/dev/null || true)
        case "$v" in
            ''|*[!0-9]*) ;;
            *) interval="$v" ;;
        esac
    fi
    if [ "$interval" -lt 60 ]; then
        # Polling faster than once a minute would burn rate limit and
        # GitHub's secondary-rate-limit thresholds for not much win.
        interval=60
    fi

    start_ts=$(date +%s)
    flock -n "$LOCK" s6-setuidgid alice /usr/local/bin/alice-gh-watcher \
        2>>/state/worker/gh-watcher-stderr.log || true
    end_ts=$(date +%s)
    elapsed=$(( end_ts - start_ts ))
    remainder=$(( interval - elapsed ))
    if [ "$remainder" -gt 0 ]; then
        sleep "$remainder"
    fi
done
